# templates/usr/lib/systemd/system/nexus_exporter.service.erb
[Unit]
Description=Nexus Metrics Exporter Service
Requires=nexus.service
After=nexus.service

[Service]
Type=idle
ExecStartPre=/bin/sleep 80
ExecStart=/usr/bin/python <%= @script_dir %>/exporter.py <%= @script_dir %>/config.yml

[Install]
WantedBy=multi-user.target