# templates/usr/lib/systemd/system/nexus_exporter.service.erb
[Unit]
Description=Nexus Metrics Exporter Service
Requires=nexus.service
After=nexus.service

[Service]
Type=idle
ExecStartPre=/bin/sleep 80
ExecStart=/usr/bin/python <%= @script_location %>'/exporter.py' <%= @script_location %>'/config.yml' # /opt/nexus_metrics_exporter/exporter.py /opt/nexus_metrics_exporter/config.yml

[Install]
WantedBy=multi-user.target